# litepanel/litepanel-admin.service
[Unit]
Description=Lite Hosting Panel — Admin Interface (port 2087)
After=network.target

[Service]
Type=simple
User=litepanel
Group=litepanel
WorkingDirectory=/usr/local/litepanel
EnvironmentFile=/usr/local/litepanel/.env
ExecStart=/usr/local/litepanel/venv/bin/gunicorn \
    --workers 2 \
    --bind 127.0.0.1:2087 \
    --certfile /usr/local/litepanel/ssl/panel.crt \
    --keyfile  /usr/local/litepanel/ssl/panel.key \
    --access-logfile /var/log/litepanel/admin-access.log \
    --error-logfile  /var/log/litepanel/admin-error.log \
    --timeout 60 \
    litepanel.wsgi:application
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
